# xBoxScript/crontab/add_run_xBox.sh
#!/bin/bash

echo "* * * * * sh '/data/xBoxScript/crontab/run_xBox.sh'" >>/var/spool/cron/root
service crond reload
